# distribution/osbuild-remote-worker.socket
[Unit]
Description=OSBuild remote worker API socket

[Socket]
Service=osbuild-composer.service
ListenStream=8700

[Install]
WantedBy=sockets.target
