# distribution/osbuild-composer-api.socket
[Unit]
Description=OSBuild Composer API socket

[Socket]
Service=osbuild-composer.service
ListenStream=443

[Install]
WantedBy=sockets.target
